# modules/protectmymail_mailserver/templates/exporter.service.erb
[Unit]
Description=<%= @name %>

[Service]
User=root
EnvironmentFile=/etc/sysconfig/<%= @name %>
ExecStart=/usr/local/bin/<%= @name %>

[Install]
WantedBy=multi-user.target
